# power-led-off.service
[Unit]
Description=Turn off power LED on shutdown or sleep
Before=shutdown.target sleep.target

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'sleep 1.8; /usr/local/bin/ectool led power off'

[Install]
WantedBy=sysinit.target
